# docker_transmission/crontabs/transmission_settings
# min	hour	day	month	weekday	command
*/15	*	*	*	*	/transmission_settings.py /config/settings.json